# functions
#!/bin/bash

# Import MOS functions
. /usr/local/lib/functions/mos-common

# Set variables for plugin
PLG_NAME="mcp-server"
DISPLAY_NAME="MCP Server"
PLG_DIR="/boot/optional/plugins/$PLG_NAME"
MCP_INSTALL_DIR="/opt/mos-mcp-server"
MCP_SERVICE="mos-mcp-server"
INIT_SCRIPT="/etc/init.d/${MCP_SERVICE}"
PIDFILE="/var/run/${MCP_SERVICE}.pid"
LOGFILE="/var/log/${MCP_SERVICE}.log"

# Install function - called after .deb installation
install() {
  echo "Installing MCP Server..."
  
  # Create installation directory
  mkdir -p "$MCP_INSTALL_DIR"
  
  # Copy MCP server files from plugin directory
  if [ -d "$PLG_DIR/staticfiles" ]; then
    cp -r "$PLG_DIR/staticfiles/"* "$MCP_INSTALL_DIR/"
  fi
  
  # Install npm dependencies if package.json exists
  if [ -f "$MCP_INSTALL_DIR/package.json" ]; then
    cd "$MCP_INSTALL_DIR"
    npm install --production 2>/dev/null || true
  fi
  
  # Create init.d script (sysvinit for Devuan)
  cat > "$INIT_SCRIPT" << 'EOF'
#!/bin/sh
### BEGIN INIT INFO
# Provides:          mos-mcp-server
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: MOS MCP Server
# Description:       Model Context Protocol server for MOS
### END INIT INFO

NAME="mos-mcp-server"
DAEMON="/usr/bin/node"
DAEMON_ARGS="/opt/mos-mcp-server/src/index.js"
PIDFILE="/var/run/${NAME}.pid"
LOGFILE="/var/log/${NAME}.log"

export NODE_ENV=production
export MOS_API_URL=http://127.0.0.1:3000

do_start() {
    echo "Starting $NAME..."
    if [ -f "$PIDFILE" ] && kill -0 $(cat "$PIDFILE") 2>/dev/null; then
        echo "$NAME is already running"
        return 1
    fi
    $DAEMON $DAEMON_ARGS >> "$LOGFILE" 2>&1 &
    echo $! > "$PIDFILE"
    echo "$NAME started"
}

do_stop() {
    echo "Stopping $NAME..."
    if [ ! -f "$PIDFILE" ]; then
        echo "$NAME is not running"
        return 1
    fi
    kill $(cat "$PIDFILE") 2>/dev/null
    rm -f "$PIDFILE"
    echo "$NAME stopped"
}

do_status() {
    if [ -f "$PIDFILE" ] && kill -0 $(cat "$PIDFILE") 2>/dev/null; then
        echo "$NAME is running (PID: $(cat $PIDFILE))"
    else
        echo "$NAME is not running"
    fi
}

case "$1" in
    start)   do_start ;;
    stop)    do_stop ;;
    restart) do_stop; do_start ;;
    status)  do_status ;;
    *)       echo "Usage: $0 {start|stop|restart|status}"; exit 1 ;;
esac
EOF

  chmod 755 "$INIT_SCRIPT"
  
  # Enable and start service
  update-rc.d ${MCP_SERVICE} defaults 2>/dev/null || true
  ${INIT_SCRIPT} start
  
  echo "MCP Server installed successfully"
}

# Uninstall function - called before package removal
uninstall() {
  echo "Uninstalling MCP Server..."
  
  # Stop service
  ${INIT_SCRIPT} stop 2>/dev/null || true
  
  # Remove from runlevels
  update-rc.d -f ${MCP_SERVICE} remove 2>/dev/null || true
  
  # Remove init script
  rm -f "$INIT_SCRIPT"
  
  # Remove installation directory
  rm -rf "$MCP_INSTALL_DIR"
  
  # Remove pid and log files
  rm -f "$PIDFILE" "$LOGFILE"
  
  echo "MCP Server uninstalled successfully"
}

# Start function
start() {
  ${INIT_SCRIPT} start
}

# Stop function
stop() {
  ${INIT_SCRIPT} stop
}

# Restart function
restart() {
  ${INIT_SCRIPT} restart
}

# Status function
status() {
  ${INIT_SCRIPT} status
}
